# Billing-2.0.7/scripts/logf.cron
#!/usr/local/bin/bash
base=localhost:/usr/local/var/internet.gdb
user=INTERNET
pass=ipfwlog
network=192.168.1
start_rule=2000
iip=192.168.1.1

/sbin/ipfw show | /usr/bin/grep count | /usr/local/sbin/logf ${base} ${user} ${pass} ${iip}
/sbin/ipfw zero
/usr/local/sbin/logf.rules ${base} ${user} ${pass} ${start_rule} ${network} ${iip}
